# etc/systemd/system/cmon-proxy.service.in
[Unit]
Description=Severalnines CMON Proxy
After=syslog.target
After=network.target

[Service]
ExecStart=@CCMGR_PATH@ --basedir=/usr/share/ccmgr/
PIDFile=/run/cmon-proxy.pid
Restart=always
User=s9s_cc
Group=severalnines
EnvironmentFile=-/etc/default/cmon-proxy.env
ExecReload=/bin/kill -HUP $MAINPID
WorkingDirectory=/usr/bin
LimitNOFILE=4096

[Install]
WantedBy=multi-user.target